# etc/rc.local
#!/bin/sh -e
#
# rc.local
#
#echo 1000 > /sys/module/xeno_nucleus/parameters/xenomai_gid
/opt/www/django_datarecorder/manage.py runfcgi method=threaded host=127.0.0.1 port=3033


# Get sure that sudoers has 0 as uid 
chown root:root /etc/sudoers

# Disable Standby for debugging messages
sh -c 'setterm -blank 0 -powersave off -powerdown 0 < /dev/console > /dev/console 2>&1'

sh /etc/voice_startup.sh &

exit 0
